# aios-web.service
[Unit]
Description=AI OS Agent Web Interface
After=network.target

[Service]
Type=simple
User=aryan
WorkingDirectory=/home/aryan/Doings/aiOS/aiOS
ExecStart=/usr/bin/python3 /home/aryan/Doings/aiOS/aiOS/web_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
